# supervisord.conf
[supervisord]
nodaemon=false
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/tmp/supervisord.pid

[program:sidekiq]
command=/home/app/bin/bundle exec sidekiq